# pushmute.service
[Unit]
Description=PushMute (selective mic router)
# Order after the audio stack so the loopback has PipeWire to attach to.
After=pipewire.service wireplumber.service
Wants=pipewire.service
# NOTE: WantedBy=default.target (below), NOT graphical-session.target. The latter
# is only activated by full desktop sessions (GNOME/KDE) or compositors launched
# via a systemd session manager (e.g. uwsm); a bare Hyprland/sway exec-once
# session never reaches it, so the unit would stay "enabled" yet never autostart.
# default.target is reached by every systemd --user instance. The tray tolerates
# the StatusNotifierWatcher appearing later (see src/tray.rs), so it no longer
# needs to be ordered after the graphical session to get its icon.

[Service]
Type=simple
ExecStart=%h/.local/bin/pushmute run
# Best-effort restore of the default source if the daemon is stopped/killed.
ExecStopPost=%h/.local/bin/pushmute restore
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
