# pushmute.service
[Unit]
Description=PushMute (selective mic router)
# Order after the audio stack and the graphical session so the tray (SNI) has a
# StatusNotifierWatcher to register against on cold boot.
After=pipewire.service wireplumber.service graphical-session.target
Wants=pipewire.service
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/pushmute run
# Best-effort restore of the default source if the daemon is stopped/killed.
ExecStopPost=%h/.local/bin/pushmute restore
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
